# scripts/com.travelplanner.dbbackup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.travelplanner.dbbackup</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alexkaplinsky/Desktop/Dev site/travel-planner-v2/scripts/db-backup.sh</string>
    </array>
    
    <key>StartInterval</key>
    <integer>600</integer>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/alexkaplinsky/Desktop/Dev site/travel-planner-v2/scripts/db-backup.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/alexkaplinsky/Desktop/Dev site/travel-planner-v2/scripts/db-backup-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
